# bin/robot.service
[Unit]
Description=Robot
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/projects/robot
User=pi
Group=pi
Restart=on-failure

[Install]
WantedBy=multi-user.target
